# packaging/macos/launchagent/com.langguard.cursor-otel-hook-setup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.langguard.cursor-otel-hook-setup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Library/Application Support/CursorOtelHook/setup-user.sh</string>
    </array>

    <!-- Run at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Run daily at 9am to ensure config stays synced with system -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Don't keep running, just run once -->
    <key>KeepAlive</key>
    <false/>

    <!-- Log output for debugging -->
    <key>StandardOutPath</key>
    <string>/tmp/cursor-otel-hook-setup.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cursor-otel-hook-setup.err</string>

    <!-- Nice priority - don't impact login performance -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Low priority I/O -->
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
